# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log

[program:httpd]
command=/usr/sbin/httpd -DFOREGROUND

[program:sshd]
command=/usr/sbin/sshd -D
